# launchd/com.catchingapril.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  CatchingApril web server — Flask app on :5000.

  Install via launchd/install.sh, which substitutes __USER_HOME__ with $HOME
  and copies the result to ~/Library/LaunchAgents/.

  Logs land in ~/Library/Logs/catchingapril-server.log.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.catchingapril.server</string>

  <!--
    Invoked via run_server.sh (not python directly) so the wrapper can
    source .env. Keeps secrets like GOOGLE_MAPS_KEY in .env instead of
    embedding them in this plist after install.
  -->
  <key>ProgramArguments</key>
  <array>
    <string>__USER_HOME__/dev/catching-april/run_server.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__USER_HOME__/dev/catching-april</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>__USER_HOME__/Library/Logs/catchingapril-server.log</string>
  <key>StandardErrorPath</key>
  <string>__USER_HOME__/Library/Logs/catchingapril-server.log</string>
</dict>
</plist>
